# utils/package-overlay.d/rpm/lib/systemd/system/travis-worker.service
[Unit]
Description=Travis Worker

[Service]
Type=simple
User=travis
Group=travis
ExecStartPre=/bin/mkdir -p /var/tmp/travis-run.d
ExecStartPre=/bin/cp /usr/local/bin/travis-worker /var/tmp/travis-run.d/%n
ExecStart=/bin/sh -c 'for config_file in travis-enterprise %n %n-local; do if [ -f /etc/default/$config_file ]; then source /etc/default/$config_file; fi; done; export GOMAXPROCS=$(nproc); exec /var/tmp/travis-run.d/%n'
ExecStopPost=/bin/sleep 5
Restart=always

[Install]
WantedBy=multi-user.target

# vim:filetype=systemd
